# config/templates/purge_varnish.cron.template
# ban varnish and delete log file if exists
* *    * * *   root    ${buildout:directory}/bin/purge_varnish.sh > ${buildout:directory}/var/log/last_purge_run.log 2>&1
